# minecraft/scripts/fishcraft.cron
@reboot sh -c "cd /home/minecraft && /usr/bin/java -Xmx2048M -Xms1024M -jar /home/minecraft/server.jar nogui"
